# systemd/mysql-db-backup.service
[Unit]
Description=Running mysql-db-backup

[Service]
ExecStart=mysql-db-backup

User=backup
Group=backup