# systemd/mysql-db-backup.timer
[Unit]
Description=Running mysql-db-backup

[Install]
WantedBy=timers.target